# aur/pacman-sync-client@.service
[Unit]
Description=Pacman Sync Utility Client (Instance %i)
Documentation=man:pacman-sync-client(1)
After=graphical-session.target
Wants=graphical-session.target
PartOf=graphical-session.target
Requisite=graphical-session.target

[Service]
Type=exec
ExecStart=/usr/bin/pacman-sync-client --instance=%i
Restart=always
RestartSec=5
RestartPreventExitStatus=0
TimeoutStartSec=30
TimeoutStopSec=10
KillMode=mixed
KillSignal=SIGTERM

# Environment variables for GUI applications
Environment=DISPLAY=:%i
Environment=QT_QPA_PLATFORM=xcb
Environment=XDG_RUNTIME_DIR=/run/user/%i
Environment=WAYLAND_DISPLAY=wayland-%i
Environment=QT_WAYLAND_DISABLE_WINDOWDECORATION=1
Environment=QT_AUTO_SCREEN_SCALE_FACTOR=1
Environment=QT_SCALE_FACTOR_ROUNDING_POLICY=RoundPreferFloor

# Instance-specific configuration
Environment=PACMAN_SYNC_INSTANCE=%i
Environment=PACMAN_SYNC_CONFIG_DIR=%h/.config/pacman-sync-utility/%i

# Security settings (appropriate for user service)
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RemoveIPC=yes

# Resource limits
MemoryMax=512M
TasksMax=100

[Install]
WantedBy=default.target